# axon-api-proxy/axon-api-proxy.service
[Unit]
Description=Axon API Proxy - OpenAI-compatible API aggregator
After=network.target

[Service]
Type=simple
User=ramlej
WorkingDirectory=/home/ramlej/axon-api-proxy
ExecStart=/home/ramlej/axon-api-proxy/axon-api-proxy
EnvironmentFile=-/home/ramlej/axon-api-proxy/.env
Environment=RUST_LOG=axon_api_proxy=info
Environment=MODELS_TOML=/home/ramlej/axon-api-proxy/models.toml
Environment=ENV_FILE=/home/ramlej/axon-api-proxy/.env
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
